# resources/rootfs/etc/systemd/system/chg_aptlist.service
[Unit]
Description=cheng apt sources.list

[Service]
Type=oneshot
ExecStart=/usr/sbin/chg_aptlist

[Install]
WantedBy=graphical.target
